# .config/autostart/blueman-manager.desktop
[Desktop Entry]
Name=Bluetooth Manager
Type=Application
Icon=blueman
Exec=blueman-applet
StartupNotify=false
