# .config/autostart/light-locker.desktop
[Desktop Entry]
Name=Lock Screen
Type=Application
Exec=light-locker
StartupNotify=false
