# .config/autostart/mpnotd.desktop
[Desktop Entry]
Name=Music Notifications
Type=Application
Exec=mpnotd
StartupNotify=false
