# .config/autostart/numlockx.desktop
[Desktop Entry]
Name=Number Lock On
Type=Application
Exec=numlockx
StartupNotify=false
