# .config/autostart/picom.desktop
[Desktop Entry]
Name=X Compositor
Type=Application
Exec=picom
StartupNotify=false
