# .config/autostart/pulseaudio.desktop
[Desktop Entry]
Name=Sound Server
Type=Application
Exec=start-pulseaudio-x11
StartupNotify=false
