# .config/autostart/udiskie.desktop
[Desktop Entry]
Name=Removable Media
Type=Application
Exec=udiskie
StartupNotify=false
