# .config/autostart/unclutter.desktop
[Desktop Entry]
Name=Hide Mouse Cursor
Type=Application
Exec=unclutter --exclude-root
StartupNotify=false
