# .config/autostart/update-notify.desktop
[Desktop Entry]
Name=Software Updates
Type=Application
Exec=software-update-notify.sh
StartupNotify=false
